# shell/install_dev/data/etc/container/php-5.3.29/crontab
*/1 * * * * root echo "Hello world" >> /var/log/cron.log 2>&1
*/1 * * * *  root /usr/bin/php   /opt/htdocs/script/crontab/crontab.php
*/1 * * * *  root /opt/htdocs/script/queue/queue.sh /usr/bin/php
*/1 * * * *  root echo "end" >> /var/log/cron.log 2>&1





